# deploy/launchd/io.tabura.lmstudio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.tabura.lmstudio</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@LMSTUDIO_SESSION_SCRIPT@@</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tabura-lmstudio.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tabura-lmstudio.log</string>
</dict>
</plist>
