# deploy/launchd/io.tabura.codex-llm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.tabura.codex-llm</string>
  <key>ProgramArguments</key>
  <array>
    <string>@@LLM_SETUP_SCRIPT@@</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>TABURA_LLM_PRESET</key>
    <string>codex-gpt-oss-120b</string>
    <key>LLAMA_SERVER_BIN</key>
    <string>@@LLAMA_SERVER_BIN@@</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/tabura-codex-llm.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/tabura-codex-llm.log</string>
</dict>
</plist>
