# deploy/systemd/user/tabura-mcp.service
[Unit]
Description=Tabura Local MCP Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabura
PassEnvironment=DISPLAY WAYLAND_DISPLAY XAUTHORITY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS
Environment=HOME=%h
Environment=SHELL=/usr/bin/bash
ExecStart=/usr/bin/bash -lc 'source ~/.bashrc >/dev/null 2>&1 || true; cd %h/code/assi/tabura; exec /usr/bin/env go run ./cmd/tabura serve --project-dir %h/code/assi/tabura --host 127.0.0.1 --port 9420'
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

[Install]
WantedBy=default.target
