# deploy/systemd/user/tabura-codex-llm.service
[Unit]
Description=Tabura Local Codex LLM (gpt-oss-120b via llama.cpp)
After=network.target

[Service]
Type=simple
WorkingDirectory=@@REPO_ROOT@@
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TABURA_LLM_PRESET=codex-gpt-oss-120b
Environment=LLAMA_SERVER_BIN=@@LLAMA_SERVER_BIN@@
ExecStart=@@REPO_ROOT@@/scripts/setup-local-llm.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=default.target
