# deploy/systemd/user/tabula-ptyd.service
[Unit]
Description=Tabula Local PTY Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabula
Environment=HOME=%h
Environment=SHELL=/usr/bin/bash
ExecStart=/usr/bin/bash -lc 'source ~/.bashrc >/dev/null 2>&1 || true; cd %h/code/assi/tabula; exec /usr/bin/env go run ./cmd/tabula ptyd --data-dir %h/.local/share/tabula-ptyd --host 127.0.0.1 --port 9333'
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
