# deploy/systemd/user/tabura-stt.service
[Unit]
Description=Tabura voxtype STT Sidecar (OpenAI-compatible)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabula
Environment=HOME=%h
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TABURA_STT_LANGUAGE=de,en
Environment=TABURA_STT_MODEL=large-v3-turbo
ExecStart=%h/code/assi/tabula/scripts/setup-whisper-stt.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=15

[Install]
WantedBy=default.target
