# deploy/systemd/user/tabula-web.service
[Unit]
Description=Tabula Web UI
After=network.target tabula-mcp.service tabula-codex-app-server.service tabula-voxtype-mcp.service
Wants=tabula-mcp.service tabula-codex-app-server.service tabula-voxtype-mcp.service

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabula
PassEnvironment=DISPLAY WAYLAND_DISPLAY XAUTHORITY XDG_RUNTIME_DIR DBUS_SESSION_BUS_ADDRESS
Environment=HOME=%h
Environment=SHELL=/usr/bin/bash
Environment=TABULA_NO_AUTH=1
ExecStart=/usr/bin/bash -lc 'source ~/.bashrc >/dev/null 2>&1 || true; cd %h/code/assi/tabula; exec /usr/bin/env go run ./cmd/tabula web --project-dir %h/code/assi/tabula --data-dir %h/.local/share/tabula-web --host 0.0.0.0 --port 8420 --local-mcp-url http://127.0.0.1:9420/mcp --app-server-url ws://127.0.0.1:8787 --dev-runtime'
Restart=on-failure
RestartSec=2
TimeoutStopSec=5

[Install]
WantedBy=default.target
