# deploy/systemd/user/tabura-dev-watch.service
[Unit]
Description=Watch Tabura source for dev reload
After=tabura-mcp.service tabura-codex-app-server.service tabura-web.service
Wants=tabura-mcp.service tabura-codex-app-server.service tabura-web.service

[Service]
Type=simple
WorkingDirectory=%h/code/assi/tabura
ExecStart=%h/code/assi/tabura/scripts/tabura-dev-watch.sh
Restart=always
RestartSec=1

[Install]
WantedBy=default.target
